# plists/com.bootindustries.mlx-lm-41961.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.bootindustries.mlx-lm-41961</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/mlx_lm.server</string>
    <string>--model</string>
    <string>/Users/nesbitt/models/Qwen3.5-4B-MLX-8bit</string>
    <string>--host</string>
    <string>0.0.0.0</string>
    <string>--port</string>
    <string>41961</string>
    <string>--use-default-chat-template</string>
  </array>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>15</integer>
  <key>StandardOutPath</key>
  <string>/Users/nesbitt/Library/Logs/factory/mlx-41961.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/nesbitt/Library/Logs/factory/mlx-41961.log</string>
</dict>
</plist>
