# plists/com.bootindustries.gsd-sidecar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.bootindustries.gsd-sidecar</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/python3</string>
    <string>/Users/nesbitt/dev/factory/portal/server.py</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>GSD_HOST</key>
    <string>127.0.0.1</string>
    <key>GSD_PORT</key>
    <string>41911</string>
    <key>GSD_DATA_ROOT</key>
    <string>/Users/nesbitt/dev/factory/portal</string>
  </dict>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>10</integer>
  <key>StandardOutPath</key>
  <string>/Users/nesbitt/Library/Logs/factory/gsd-sidecar.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/nesbitt/Library/Logs/factory/gsd-sidecar.log</string>
</dict>
</plist>
