# agents/ig88/scripts/matrix-coordinator.service
[Unit]
Description=IG-88 Matrix Coordinator
Documentation=https://github.com/chrislyons/factory
After=network.target docker.service
Wants=network.target

[Service]
Type=simple
User=nesbitt
Group=nesbitt
WorkingDirectory=/home/nesbitt/projects/ig88/src

# Environment
Environment=NODE_ENV=production
Environment=HOME=/home/nesbitt
EnvironmentFile=/home/nesbitt/.config/ig88/coordinator.env

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/nesbitt/projects/ig88/logs /home/nesbitt/.config/ig88/checkpoints
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

# Execution
ExecStart=/usr/bin/node /home/nesbitt/projects/ig88/src/dist/matrix-coordinator.js
Restart=always
RestartSec=10

# Logging
StandardOutput=append:/home/nesbitt/projects/ig88/logs/coordinator.log
StandardError=append:/home/nesbitt/projects/ig88/logs/coordinator-error.log

[Install]
WantedBy=multi-user.target
