# agents/ig88/scripts/services/blackbox.target
[Unit]
Description=Blackbox Infrastructure Stack
Documentation=file:///home/nesbitt/projects/ig88/docs/startup.md
After=network-online.target docker.service tailscaled.service
Requires=docker.service tailscaled.service
Wants=network-online.target
Wants=qdrant.service
Wants=ollama.service
Wants=graphiti.service
Wants=pantalaimon.service
Wants=matrix-coordinator.service
Wants=heartbeat.service

[Install]
WantedBy=multi-user.target
