# agents/ig88/scripts/services/graphiti.service
[Unit]
Description=Graphiti Knowledge Graph (FalkorDB + MCP Server)
Documentation=https://github.com/getzep/graphiti
After=network-online.target docker.service ollama.service
Requires=docker.service
Wants=network-online.target ollama.service

[Service]
Type=oneshot
RemainAfterExit=yes
User=nesbitt
Group=docker
WorkingDirectory=/home/nesbitt/projects/graphiti

# Start docker compose stack
ExecStart=/usr/bin/docker compose up -d
ExecStop=/usr/bin/docker compose down

# Health check - wait for FalkorDB to be ready
ExecStartPost=/bin/bash -c 'for i in {1..60}; do docker compose exec -T graphiti-falkordb redis-cli ping && exit 0; sleep 2; done; echo "FalkorDB not ready"; exit 1'

# Restart policy handled by docker compose restart: always
Restart=on-failure
RestartSec=30

[Install]
WantedBy=blackbox.target
