# agents/ig88/scripts/services/heartbeat.service
[Unit]
Description=Blackbox Heartbeat Monitor
Documentation=file:///home/nesbitt/projects/ig88/docs/heartbeat.md
After=network-online.target matrix-coordinator.service
Wants=network-online.target

[Service]
Type=simple
User=nesbitt
Group=nesbitt
WorkingDirectory=/home/nesbitt/projects/ig88/src

# Environment
Environment=NODE_ENV=production
Environment=HOME=/home/nesbitt
Environment=MATRIX_ROOM_ID=!MDVmYJtAiHZoBfaQdK:matrix.org
Environment=MATRIX_TOKEN_FILE=/home/nesbitt/.config/ig88/matrix_token_boot
Environment=MATRIX_HOMESERVER=https://matrix.org

# Main process
ExecStart=/usr/bin/node /home/nesbitt/projects/ig88/src/dist/heartbeat.js

# Restart policy - always restart, it's our watchdog
Restart=always
RestartSec=30
StartLimitIntervalSec=300
StartLimitBurst=10

# Logging
StandardOutput=append:/home/nesbitt/projects/ig88/logs/heartbeat.log
StandardError=append:/home/nesbitt/projects/ig88/logs/heartbeat-error.log

[Install]
WantedBy=blackbox.target
