# agents/ig88/scripts/services/matrix-coordinator.service
[Unit]
Description=IG-88 Matrix Coordinator (Multi-Agent)
Documentation=file:///home/nesbitt/projects/ig88/docs/bkx/BKX015
After=network-online.target tailscaled.service qdrant.service ollama.service pantalaimon.service
Wants=network-online.target
Requires=tailscaled.service
Wants=qdrant.service ollama.service pantalaimon.service

[Service]
Type=simple
User=nesbitt
Group=nesbitt
WorkingDirectory=/home/nesbitt/projects/ig88/src

# Environment
Environment=NODE_ENV=production
Environment=HOME=/home/nesbitt
Environment=CONFIG_PATH=/home/nesbitt/projects/ig88/config/agent-config.yaml
EnvironmentFile=-/home/nesbitt/.config/ig88/coordinator.env

# Wait for Tailscale to be connected before starting
ExecStartPre=/bin/bash -c 'for i in {1..60}; do tailscale status --json | grep -q "Online" && exit 0; sleep 2; done; echo "Tailscale not connected"; exit 1'

# Main process
ExecStart=/usr/bin/node /home/nesbitt/projects/ig88/src/dist/matrix-coordinator.js

# Restart policy
Restart=always
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/nesbitt/projects/ig88/logs
ReadWritePaths=/home/nesbitt/.config/ig88
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

# Logging
StandardOutput=append:/home/nesbitt/projects/ig88/logs/coordinator.log
StandardError=append:/home/nesbitt/projects/ig88/logs/coordinator-error.log

# Graceful shutdown
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=blackbox.target
