# agents/ig88/scripts/services/ollama.service
[Unit]
Description=Ollama LLM Server (Embeddings)
Documentation=https://ollama.ai/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nesbitt
Group=nesbitt

ExecStart=/usr/local/bin/ollama serve
Restart=always
RestartSec=10

# Environment
Environment=OLLAMA_HOST=0.0.0.0:11434
Environment=OLLAMA_MODELS=/home/nesbitt/.ollama/models

# Resource limits
MemoryMax=4G
CPUQuota=300%

# Health check
ExecStartPost=/bin/bash -c 'for i in {1..30}; do curl -s http://localhost:11434/api/tags && exit 0; sleep 1; done; exit 1'

[Install]
WantedBy=blackbox.target
